# autogauge.service
[Unit]
Description=AutoGauge dashboard display
After=multi-user.target

[Service]
Type=simple
User=alvis
Group=alvis
WorkingDirectory=/home/alvis/gauge
ExecStart=/usr/bin/python3 -u /home/alvis/gauge/main.py
Restart=on-failure
RestartSec=2
StartLimitIntervalSec=60
StartLimitBurst=10
Environment=PYTHONUNBUFFERED=1
Environment=SDL_VIDEODRIVER=offscreen

[Install]
WantedBy=multi-user.target
